# cloud-saves.service
[Unit]
Description=SillyTavern Cloud Backup Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cloud-saves
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10

# 日志
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
